# deploy/supervisor_conf.d/webapp.conf
; =======
;  Django
; =======

[program:webapp]
directory=/var/projects/mtgpl/
environment=DJANGO_SETTINGS_MODULE="config.settings"
command=gunicorn config.wsgi:application -c /var/projects/mtgpl/deploy/gunicorn.conf
autostart=true
autorestart=true
stdout_logfile = /var/log/webapp/gunicorn_supervisor.log             ; Where to write log messages
redirect_stderr=True
user=webapp
